# launchagents/com.mgc.salvage-task-poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mgc.salvage-task-poller</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jayden.csai/bin/salvage-task-poller</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/jayden.csai/claude-agent/logs/task-poller.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jayden.csai/claude-agent/logs/task-poller.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/jayden.csai/bin</string>
        <key>HOME</key>
        <string>/Users/jayden.csai</string>
    </dict>
</dict>
</plist>
