# launchagents/com.salvage.cron.linkedin-jayden-context.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.salvage.cron.linkedin-jayden-context</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jayden.csai/bin/salvage</string>
        <string>--task</string>
        <string>general</string>
        <string>--workspace</string>
        <string>/Users/jayden.csai/clawd</string>
        <string>--fresh</string>
        <string>Run linkedin-content skill Mode 1: Review daily notes. Identify 1-2 LinkedIn-worthy moments from Jayden's work this week. Draft posts.</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>9</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/jayden.csai/claude-agent/logs/cron-linkedin-jayden-context.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jayden.csai/claude-agent/logs/cron-linkedin-jayden-context.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/jayden.csai/bin</string>
        <key>HOME</key>
        <string>/Users/jayden.csai</string>
    </dict>
</dict>
</plist>