# launchagents/com.mgc.whisper-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mgc.whisper-proxy</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/jayden.csai/claude-agent/whisper-proxy.js</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/jayden.csai/claude-agent/logs/whisper-proxy.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jayden.csai/claude-agent/logs/whisper-proxy.log</string>
</dict>
</plist>
